# config/systemd/user/openrgb.service
[Unit]
Description=OpenRGB - RGB Lighting Control
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/flatpak run org.openrgb.OpenRGB --server -p Ideal
Restart=on-failure
RestartSec=5
SuccessExitStatus=0
# If running under Wayland + server mode, DISPLAY is typically not needed.
# Uncomment if required by your environment:
# Environment=DISPLAY=:0

[Install]
WantedBy=default.target
